# crontab
00 00 * * * /usr/bin/duply backupconfiguration backup --allow-source-mismatch >> /var/log/cron.log 2>&1

